# deploy/systemd/gofiber-api.service
[Unit]
Description=GoFiber API Service
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=gofiber
Group=gofiber
WorkingDirectory=/home/gofiber/serkk-backend
ExecStart=/home/gofiber/serkk-backend/bin/api
Restart=always
RestartSec=10

# Environment file
EnvironmentFile=/home/gofiber/serkk-backend/.env

# Logging
StandardOutput=append:/var/log/gofiber-api/access.log
StandardError=append:/var/log/gofiber-api/error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/gofiber/serkk-backend

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
